# roles/redis-commander/files/redis-commander.service
[Unit]
After=redis-server.service

[Service]
User=redis
ExecStart=/usr/local/bin/redis-commander
Restart=on-failure
